# meshtastic-bbs/mesh-bbs-healthcheck.service
# /etc/systemd/system
[Unit]
Description=Health check for mesh-bbs container
After=mesh-bbs.service

[Service]
Type=oneshot
ExecStart=/bin/bash -c "\
  HEALTH=$(/usr/bin/docker inspect --format='{{.State.Health.Status}}' mesh-bbs 2>/dev/null || echo 'missing'); \
  case \"$HEALTH\" in \
    unhealthy) \
      logger 'mesh-bbs healthcheck: Container unhealthy, restarting service'; \
      systemctl restart mesh-bbs.service || logger 'mesh-bbs healthcheck: Failed to restart service'; \
      ;; \
    missing) \
      if systemctl is-active --quiet mesh-bbs.service; then \
        logger 'mesh-bbs healthcheck: Container missing but service active, restarting'; \
        systemctl restart mesh-bbs.service || logger 'mesh-bbs healthcheck: Failed to restart service'; \
      fi; \
      ;; \
    healthy) \
      logger 'mesh-bbs healthcheck: Container healthy'; \
      ;; \
  esac"

